# recipes-connectivity/libpcap/libpcap/rpcapd.service
[Unit]
Description=Rpcapd Service


[Service]
ExecStart=/usr/sbin/rpcapd -n
Restart=always
RestartSec=2s

[Install]
WantedBy=multi-user.target